# ex04/00-usb-keyboards.rules
SUBSYSTEM=="input", SUBSYSTEMS=="usb", ACTION=="add", RUN+="/sbin/insmod /root/LP/main.ko", OWNER="root"
